# backend/crontab
0 */6 * * * root /usr/local/bin/python -m app.jobs.collect_news 2>&1 | tee -a /var/log/cron.log

# O cron exige uma linha em branco no final do arquivo para funcionar corretamente.
